# deploy/buildrent.service
[Unit]
Description=BuildRent Spring Boot Application
After=network.target mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/buildrent
ExecStart=/usr/bin/java -jar /opt/buildrent/BuildRent-0.0.1-SNAPSHOT.jar --spring.profiles.active=prod
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=DB_PASSWORD=192856qwaS
Environment=JAVA_OPTS=-Xms256m -Xmx512m

[Install]
WantedBy=multi-user.target
